# rootfs/etc/services.d/media-services/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Stop the media services 
# ==============================================================================

cd /etc
docker-compose down